# config/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true                ; (start in foreground if true;default false)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:gateway]
command=/opt/klaczng/gateway/bin/gateway --gateway-rpc-endpoint tcp://*:9001 --gateway-pub-endpoint tcp://*:2137 --ephemeral-rpc-endpoint tcp://localhost:3543 --irc-config=/srv/klaczng/config/gateway.json
user=klaczng
startsecs=5
startretries=5
autorestart=unexpected
stdout_logfile=/tmp/gateway_stdout.log
stderr_logfile=/tmp/gateway_stderr.log


[program:ephemeral]
command=/opt/klaczng/ephemeral/bin/ephemeral --rpc-endpoint tcp://*:3543
user=klaczng
startsecs=5
startretries=5
autorestart=unexpected
stdout_logfile=/tmp/ephemeral_stdout.log
stderr_logfile=/tmp/ephemeral_stderr.log


[program:stateless]
command=/opt/klaczng/stateless/bin/stateless --gateway-pub-endpoint tcp://localhost:2137 --gateway-rpc-endpoint tcp://localhost:9001 --ephemeral-pub-endpoint tcp://localhost:3543
user=klaczng
startsecs=5
startretries=5
autorestart=unexpected
stdout_logfile=/tmp/stateless_stdout.log
stderr_logfile=/tmp/stateless_stderr.log


